# bin/cmder/config/profile.d/laragon.cmd
@echo off
call %CMDER_ROOT%\..\..\etc\cmder\laragon.cmd
exit /b 0
